# bootc/rootfs/etc/containers/systemd/users/1000/service-gator.container
[Unit]
Description=service-gator scoped service MCP server (rootless Podman)
After=network-online.target

[Container]
Image=ghcr.io/cgwalters/service-gator:placeholder
ContainerName=service-gator
Pull=missing
RunInit=true
# Pin SELinux MCS to the same pair as clawx.container — see the comment on
# the workspaces Volume= in clawx.container for the full rationale. Both
# values MUST match or the shared :Z workspace mount races on relabel.
SecurityLabelLevel=s0:c200,c500
Volume=%h/.config/service-gator:/etc/service-gator:Z
Volume=%h/workspaces:/workspaces:Z
Network=clawx-isolated.network
Environment=GH_TOKEN_FILE=/run/secrets/gh_token
Environment=GITLAB_TOKEN_FILE=/run/secrets/gitlab_token
Environment=FORGEJO_TOKEN_FILE=/run/secrets/forgejo_token
Environment=JIRA_API_TOKEN_FILE=/run/secrets/jira_api_token
Exec=--mcp-server 0.0.0.0:8080 --scope-file /etc/service-gator/scopes.json

[Service]
ExecStartPre=/usr/libexec/tank-os/bootstrap-clawx
ExecStartPre=/usr/libexec/tank-os/bootstrap-service-gator
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
