# bootc/rootfs/etc/systemd/system/clawx-opencode-config.path
[Unit]
Description=Watch clawx agent.env and regenerate opencode config

[Path]
PathChanged=/var/home/clawx/.clawx/agent.env

[Install]
WantedBy=multi-user.target
